# archiso-profile/airootfs/etc/systemd/system/proton-ge-updater.service
[Unit]
Description=Proton-GE Automatic Updater
After=network-online.target
Wants=network-online.target
Documentation=https://github.com/GloriousEggroll/proton-ge-custom

[Service]
Type=oneshot
ExecStart=/usr/local/bin/proton-ge-updater.sh
StandardOutput=journal
StandardError=journal

# Security hardening
ProtectSystem=strict
ReadWritePaths=/usr/share/steam/compatibilitytools.d /home /tmp /var/log
ProtectHome=read-only
NoNewPrivileges=false
PrivateTmp=true
